# home/dot_config/systemd/user/distrobox-upgrade-fedora-dev.timer
[Unit]
Description=Daily update for fedora-dev distrobox

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1d

[Install]
WantedBy=timers.target
