# home/dot_config/systemd/user/distrobox-upgrade-fedora-dev.service
[Unit]
Description=Update fedora-dev distrobox

[Service]
Type=oneshot
ExecStart=/usr/bin/distrobox-upgrade fedora-dev
StandardOutput=journal
StandardError=journal
